# automation-pilot.service
[Unit]
Description=automation-pilot — Automation Execution Orchestrator
After=network.target

[Service]
Type=simple
User=autopilot
Group=autopilot
WorkingDirectory=/opt/automation-pilot
EnvironmentFile=/opt/automation-pilot/.env
Environment=AUTOPILOT_CONFIG=/opt/automation-pilot/backend/config.yaml
Environment=AUTOPILOT_DB_PATH=/var/lib/automation-pilot/autopilot.db
ExecStart=/opt/automation-pilot/.venv/bin/python run.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
